# uniwill-fan.service
[Unit]
Description=Uniwill Silent Fan Control
After=multi-user.target
ConditionPathExists=/sys/class/uniwill_fan/uniwill_fan

[Service]
Type=simple
ExecStart=/usr/local/bin/uniwill-fanctl
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
